# etc/systemd/system/local.service
[Unit]
Description=Local scripts

[Service]
ExecStart=./local.sh
